# server/usr/lib/systemd/system/postfix-get-group@.service
[Unit]
Description=postfix-get-group server service

[Service]
User=nobody
ExecStart=/usr/libexec/nethserver/postfix-get-group
StandardInput=socket
