# server/usr/lib/systemd/system/postfix-get-group.socket
[Unit]
Description = postfix-get-group server

[Socket]
ListenStream = 127.0.0.1:14444
Accept = yes

[Install]
WantedBy = sockets.target
