# filter/usr/lib/systemd/system/clamd@rspamd.service.d/nethserver.conf
[Service]
User=_rspamd
Group=mail
RuntimeDirectory=clamd@rspamd
RuntimeDirectoryMode=0755
