# common/etc/systemd/system/postfix.service.d/nethserver.conf
[Unit]
Wants=postfix-create-DH2048.service
