# p3scan/usr/lib/systemd/system/p3scan.service.d/nethserver.conf
[Service]
User=root
Group=mail
RestartSec=3
Restart=on-success
